# systemd/kiro-openai-proxy.service
[Unit]
Description=kiro OpenAI-compatible proxy
After=network.target acp-bridge.service

[Service]
ExecStart=node %h/.local/lib/kiro-openai-proxy/index.js
Restart=on-failure
Environment=ACP_BRIDGE_URL=http://127.0.0.1:7800
# KIRO_PROXY_KEY should be set in your environment, not here
# Add to ~/.bashrc: export KIRO_PROXY_KEY=your-random-string

[Install]
WantedBy=default.target
